# etc/systemd/display-off.service
[Unit]
Description=display_off

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=/home/fyr/.dots/bin/doff

[Install]
WantedBy=multi-user.target
