# etc/systemd/radeon.service
[Unit]
Description=Undervolt Radeon VII

[Service]
ExecStart=/home/fyr/.dots/bin/radeon undervolt

[Install]
WantedBy=multi-user.target
