# etc/systemd/smbcheck.service
[Unit]
Description=smbcheck

[Service]
Type=simple
ExecStart=/home/fyr/.dots/bin/smbmount

[Install]
WantedBy=multi-user.target
