# etc/systemd/smbcheck.timer
[Unit]
Description=smbtimer

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true
Unit=smbcheck.service

[Install]
WantedBy=default.target
